# systemd/pi-live-pipeline@.service
[Unit]
Description=Pi Live Detection Pipeline for %i
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi-live-detect-rstp
ExecStart=/usr/bin/env sh -lc "source /home/pi/pi-live-detect-rstp/.venv/bin/activate && python -m app.entrypoints.pipeline_service %i"
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
